# Raspberry Pi/LCD Examples/c/lcdip.service
[Unit]
Description=lcdip
After=network.target

[Service]
ExecStart=/home/pi/lcd_examples/lcd_print_ip
WorkingDirectory=/home/pi/lcd_examples
StandardOutput=inherit
StandardError=inherit
User=root

[Install]
WantedBy=multi-user.target